# .config/systemd/user/gc.service
[Unit]
Description= nix-cache garbage collection

[Service]
Type=oneshot
ExecStart=nix-store --gc
